# src/etc/crontabs/root
*/1 * * * * /usr/bin/logtemp /root/temp_log.csv